# webmonitor.service
[Unit]
Description=WebMonitor System Resource Monitor
After=network.target

[Service]
Type=simple
# Note: User and WorkingDirectory will be auto-configured by install_service.sh
User=your-username
WorkingDirectory=/path/to/webmonitor
ExecStart=/usr/bin/node server_linux.js
Restart=always
Environment=NODE_ENV=production PORT=3000

[Install]
WantedBy=multi-user.target